# install/runConfUser.service
[Unit]
Description=My script

[Service]
ExecStart=/archInstall/config-user.sh

[Install]
WantedBy=multi-user.target